# docker/container-files/etc/crontabs/www
* * * * * flock -nx /tmp/satisfy.cron.lock -c /data/www/bin/cron-rebuild.sh